# .gitpack/install/global-Linux/prerm
#!/bin/sh

#-------------------------------------------------------------------------------
# Copyright 2019 Dominik Salvet
# Copyright 2019 Alexander Agura
# github.com/dominiksalvet/asus-fan-control
#-------------------------------------------------------------------------------

# stop running asus-fan-control after boot
systemctl disable asus-fan-control.service
