# src/asus-fan-control.service
[Unit]
Description=ASUS laptops fan control

[Service]
Type=oneshot
ExecStart=/usr/bin/asus-fan-control

[Install]
WantedBy=multi-user.target
